# systemd/mako.service
[Unit]
Description=mako: A lightweight Wayland notification daemon

[Service]
ExecStart=mako

[Install]
WantedBy=default.target
